# service/redmaker.service
[Unit]
Description=Redundancy-Maker service
After=multi-user.target